# archive/homelab-prod/scripts/pi-power-optimize.service
[Unit]
Description=Raspberry Pi Power Optimization
After=network.target
Documentation=https://github.com/IamVanshKhanna/raspi-devops-homelab

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pi-power-optimize.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/sys /proc /boot/firmware

[Install]
WantedBy=multi-user.target