# archive/homelab-prod/scripts/rclone-replicate.service
[Unit]
Description=Restic Cross-Region Replication via rclone
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=vansh
Group=vansh
Environment=RCLONE_CONFIG=/home/vansh/.config/rclone/rclone.conf
Environment=B2_ACCOUNT_ID=%i
Environment=B2_ACCOUNT_KEY=%i
Environment=AWS_ACCESS_KEY_ID=%i
Environment=AWS_SECRET_ACCESS_KEY=%i
Environment=GOOGLE_APPLICATION_CREDENTIALS=/home/vansh/.config/gcloud/service-account.json
Environment=RCLONE_TRANSFERS=4
Environment=RCLONE_CHECKERS=8
Environment=RCLONE_RETRIES=3
Environment=RCLONE_BWLIMIT=50M
Environment=RCLONE_LOG_LEVEL=INFO
Environment=TELEGRAM_BOT_TOKEN=%i
Environment=TELEGRAM_CHAT_ID=%i
WorkingDirectory=/home/vansh/homelab-prod
ExecStart=/home/vansh/homelab-prod/scripts/replicate-restic.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rclone-replicate

[Install]
WantedBy=multi-user.target